# root/etc/s6-overlay/s6-rc.d/init-rainloop/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# rainloop php files
if [[ ! -d /config/www/rainloop ]]; then
    cp -R /rainloop /config/www/rainloop
fi
